# supervisor/confd.ini
[program:v2ray_server]
directory = /usr/local/apps/v2ray
command = /usr/local/apps/v2ray/v2ray -config /root/projects/pri/v2ray/config_server.tls.json
user = root
priority = 1
startsecs = 3
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:v2ray_client]
directory = /usr/local/apps/v2ray
command = /usr/local/apps/v2ray/v2ray -config /root/projects/pri/v2ray/config_client.tls.json
user = root
priority = 1
startsecs = 3
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:nginx]
directory = /usr/local/apps/
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf  -g 'daemon off;'
user = root
priority = 1
startsecs = 3
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:mongodb]
directory = /usr/local/apps/mongodb
command = /usr/local/apps/mongodb/bin/mongod --bind_ip 0.0.0.0 --dbpath /usr/local/apps/mongodb/data --port 27018 --auth
user = root
priority = 1
startsecs = 5
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:es]
directory = /usr/local/apps/elasticsearch-7.6.0
command = bash ./bin/elasticsearch
user = es
priority = 1
startsecs = 15
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:kibana]
directory = /usr/local/apps/kibana-7.6.0-linux-x86_64
command = bash ./bin/kibana
user = es
priority = 1
startsecs = 15
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log


[program:etcd]
directory = /usr/local/apps/etcd-v3.4.4-linux-amd64
command = /usr/local/apps/etcd-v3.4.4-linux-amd64/etcd --listen-client-urls 'http://0.0.0.0:2379' --advertise-client-urls 'http://0.0.0.0:2379'
user = root
priority = 1
startsecs = 5
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:redis]
directory = /usr/local/apps/redis
command = /usr/local/apps/redis/src/redis-server /usr/local/apps/redis/redis.conf
user = root
priority = 1
startsecs = 5
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:goproxy]
directory = /usr/local/apps/goproxy
command = /usr/local/apps/goproxy/proxy http -t tls -p ":38080" -C proxy.crt -K proxy.key 
user = root
priority = 1
startsecs = 5
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:frp]
directory = /usr/local/apps/frp
command = /usr/local/apps/frp/frps -c ./frps.ini
user = root
priority = 1
startsecs = 5
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log

[program:site]
directory = /root/projects/renyanSite
command = python manage.py runserver 0.0.0.0:80
user = root
priority = 1
startsecs = 5
startretries = 3
autostart = false
autorestart = true
redirect_stderr = true
stopasgroup=true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 2
stdout_logfile = /var/log/supervisor/%(program_name)s.log
